# ledmatrix
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Alpine Linux autostart service."

BASE_DIR=/usr/local/$RC_SVCNAME
CONFIG_FILE="/etc/${RC_SVCNAME}.ini"

command="$BASE_DIR/.venv/bin/python"
command_args="main.py --config-file $CONFIG_FILE"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
start_stop_daemon_args="--chdir $BASE_DIR"


depend() {
    need net
}

reload() {
    ebegin "Reloading $RC_SVCNAME"
    # send SIGHUP(1) to process
    start-stop-daemon --exec $command \
                      --pidfile $pidfile \
                      -s 1
    eend $?
}
