# ledmatrix
#!/sbin/openrc-run

name=$RC_SVCNAME
description="Alpine Linux autostart service."

BASE_DIR=/usr/local/$RC_SVCNAME

cfgfile="/etc/conf.d/$RC_SVCNAME"
pidfile="/run/$RC_SVCNAME.pid"

command="$BASE_DIR/.venv/bin/python"
command_args="main.py -c ${cfgfile}"
command_background="yes"
start_stop_daemon_args="--chdir $BASE_DIR"


depend() {
    use net
}
